# debian/postrm
#!/bin/sh

case "$1" in
  purge)
    update-rc.d hostapd remove

    ;;


    *)
    ;;
esac

#DEBHELPER#
